# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}
UI_PORT=${UI_PORT:-8082}

umask "$UMASK_SET"

exec s6-setuidgid abc /usr/local/bin/flexget daemon start --daemonize
